# docker/service/redis/run
#!/bin/sh

exec chpst \
   -u redis \
   -U redis \
   redis-server /etc/redis.conf
