# docker/service/celery/run
#!/bin/sh

cd /src/

export DOORMAN_ENV=prod
export DOORMAN_SETTINGS=$PWD/settings.cfg

exec chpst \
  -u doorman \
  -U doorman \
  celery \
    worker \
    -A doorman.worker:celery \
    -l INFO
